# Build/deb-meta/etc/systemd/system/mqttHub.service
[Unit]
Description=mqttHub
After=network.target

[Service]
ExecStart=/opt/mqttHub/mqttHub
WorkingDirectory=/opt/mqttHub
StandardOutput=inherit
StandardError=inherit
Restart=always
User=mqtthubsrv

[Install]
WantedBy=multi-user.target